# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10

[program:fastapi]
command=uvicorn server.app:app --host 0.0.0.0 --port 7860
directory=/app
autostart=true
autorestart=true
priority=20

[program:easy_attack]
command=/bin/bash /app/scripts/easy_attack.sh
autostart=false
autorestart=false
priority=30

[program:medium_attack]
command=/bin/bash /app/scripts/medium_attack.sh
autostart=false
autorestart=false
priority=40

[program:hard_attack]
command=/bin/bash /app/scripts/hard_attack.sh
autostart=false
autorestart=false
priority=50